# gateway/gateway-forwarder/DEBIAN/postinst
#!/bin/sh
systemctl enable lorapktfwd
systemctl start lorapktfwd
